# resources/mailpit.conf
[program:mailpit]
command=/usr/local/bin/mailpit --max 5000
autostart=true
autorestart=true
stderr_logfile=/var/log/mailpit.err.log
stdout_logfile=/var/log/mailpit.out.log
