# resources/mailhog.conf
[program:mailhog]
command=/usr/local/bin/mailhog
autostart=true
autorestart=true
stderr_logfile=/var/log/mailhog.err.log
stdout_logfile=/var/log/mailhog.out.log
